# backend/deploy/systemd/cloudflare-tunnel.service
[Unit]
Description=Cloudflare Tunnel for MVP Backend
Requires=docker.service
After=docker.service backend.service
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/mvp-backend
ExecStart=/usr/bin/docker-compose -f deploy/cloudflare/docker-compose.tunnel.yml up -d
ExecStop=/usr/bin/docker-compose -f deploy/cloudflare/docker-compose.tunnel.yml down
TimeoutStartSec=0
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
